# deploy/systemd/openclaw-tmp-reaper.service
[Unit]
Description=Reap orphaned OpenClaw render temp dirs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/openclaw-reap-tmp
