# deploy/systemd/openclaw-backup.service
[Unit]
Description=Backup OpenClaw workspace state to restic repository

[Service]
Type=oneshot
ExecStart=%h/.local/bin/openclaw-backup
Nice=10
IOSchedulingClass=idle
