# deploy/systemd/openclaw-backup.timer
[Unit]
Description=Daily backup of OpenClaw workspace state

[Timer]
OnCalendar=*-*-* 03:30
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
