# deploy/systemd/openclaw-doctor-alert.service
[Unit]
Description=Probe workspace app + gateway health, alert via ntfy on state change

[Service]
Type=oneshot
ExecStart=%h/.local/bin/openclaw-doctor-alert
