# deploy/systemd/openclaw-gateway.service.d/java-ipv4.conf
[Service]
Environment=JAVA_OPTS=-Djava.net.preferIPv4Stack=true
