# deploy/systemd/openclaw-gateway.service.d/memory.conf
# Memory guardrail for the gateway (openclaw). Observed 2026-07-09: 4.1G RSS
# climbing to 6.1G within a day, 7.6G peak over 24h, no cap. MemoryHigh
# throttles + reclaims under pressure without killing; MemoryMax is a runaway
# backstop set well above the observed peak because an OOM-kill costs a
# 4-5 minute gateway cold boot. If RSS trends toward MemoryHigh, open a leak
# investigation — this is a guardrail, not the fix.
# Apply without a restart: systemctl --user daemon-reload (cgroup values
# reapply to the running unit).
[Service]
MemoryHigh=8G
MemoryMax=12G
