# deploy/systemd/openclaw-gateway.service.d/tmpdir.conf
# Render temp off tmpfs: /tmp is tmpfs with per-user quota (usrquota); the
# gateway's rf_*/ge_* podcast render dirs (100s of MB each) filled the quota on
# 2026-07-02 and crash-looped the gateway with EDQUOT ("UNKNOWN: unknown error,
# write"). Real disk instead; reaped by openclaw-tmp-reaper.timer.
[Service]
Environment=TMPDIR=%h/.cache/openclaw-tmp
