# deploy/systemd/openclaw-gateway.service.d/voyage.conf
[Service]
EnvironmentFile=%h/.config/openclaw-secrets/gateway.env
